# admin/bootstrap/supervisor-confs/ailurusd.conf
[program:ailurusd]
autorestart=true
numprocs=2
numprocs_start=1
command=ailurusd --port=1857%(process_num)s
process_name="ailurusd-%(process_num)s"
environment=ON_AILURUS="true"
